# toyota-dashboard-user.service.example
[Unit]
Description=Toyota Dashboard Server (User)
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=%h/toyota-dashboard
Environment=HOME=%h
Environment=XDG_CONFIG_HOME=%h/.config
Environment=XDG_DATA_HOME=%h/.local/share
Environment=XDG_CACHE_HOME=%h/.cache
Environment=PYTHONPATH=%h/toyota-dashboard
ExecStart=%h/toyota-dashboard/venv/bin/python %h/toyota-dashboard/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target